# calendar/scrape.cron

#run daily at 3am 
0 3 * * * /usr/bin/wget --output-document calendar.ics https://outlook.office365.com/owa/calendar/7e6c037200a54bbb89870bdf7153444e@utoronto.ca/0222599322974e3a84fb22a53974a24b6817473586409478242/calendar.ics
